# conf/supervisor-botjagwar.conf
[program:translator_1]
command=/usr/bin/docker exec -w /dockerx -t 6feb71924bb3 /dockerx/rbcst-env/bin/python /dockerx/ctranslator.py 8888
stdout_logfile=/opt/botjagwar/user_data/translator_1_stderr.log
stderr_logfile=/opt/botjagwar/user_data/translator_1_stderr.log
autostart=true
startretries=10
user=root
stderr_logfile_maxbytes=10MB

[program:translator_2]
command=/usr/bin/docker exec -w /dockerx -t 6feb71924bb3 /dockerx/rbcst-env/bin/python /dockerx/ctranslator.py 8889
stdout_logfile=/opt/botjagwar/user_data/translator_2_stdout.log
stderr_logfile=/opt/botjagwar/user_data/translator_2_stderr.log
autostart=true
startretries=10
user=root
stderr_logfile_maxbytes=10MB




[program:load_balancer]
command=/usr/sbin/haproxy -Ws -f /opt/botjagwar/conf/haproxy.cfg -p /run/haproxy.pid -S /run/haproxy-master.sock
stdout_logfile=/opt/botjagwar/user_data/loadbalancer.log
stderr_logfile=/opt/botjagwar/user_data/loadbalancer_errors.log
autostart=true
startretries=10
user=root
stderr_logfile_maxbytes=100MB

[program:rmq_recentchanges]
command=/opt/botjagwar/pyenv/bin/python /opt/botjagwar/rabbitmq_entry_translator.py 8000
stdout_logfile=/opt/botjagwar/user_data/rmq_recentchanges_stdout.log
stderr_logfile=/opt/botjagwar/user_data/rmq_recentchanges_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=10MB

[program:dictionary_service]
command=/opt/botjagwar/pyenv/bin/python /opt/botjagwar/dictionary_service.py
stdout_logfile=/opt/botjagwar/user_data/supervisor_dictionary_service_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_dictionary_service_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB
3
[program:entry_translator_1]
command=/opt/botjagwar/pyenv/bin/python /opt/botjagwar/entry_translator.py -p 18001
stdout_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB


[program:entry_translator_2]
command=/opt/botjagwar/pyenv/bin/python /opt/botjagwar/entry_translator.py -p 18002
stdout_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:entry_translator_3]
command=/opt/botjagwar/pyenv/bin/python /opt/botjagwar/entry_translator.py -p 18003
stdout_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:entry_translator_4]
command=/opt/botjagwar/pyenv/bin/python /opt/botjagwar/entry_translator.py -p 18004
stdout_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:entry_translator_5]
command=/opt/botjagwar/pyenv/bin/python /opt/botjagwar/entry_translator.py -p 18005
stdout_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_entry_translator_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB


[program:wiktionary_irc]
command=/opt/botjagwar/pyenv/bin/python /opt/botjagwar/wiktionary_irc.py
stdout_logfile=/opt/botjagwar/user_data/supervisor_wiktionary_irc_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_wiktionary_irc_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB



[program:botjagwar_nginx]
command=/usr/sbin/nginx -g 'daemon off;' -p /opt/botjagwar-front/ -c /opt/botjagwar-front/config/nginx/nginx.conf
stdout_logfile=/opt/botjagwar/user_data/nginx_stdout.log
stderr_logfile=/opt/botjagwar/user_data/nginx_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB




[program:postgrest_2]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_2.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest2_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest2_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_3]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_3.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest3_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest3_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_4]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_4.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest4_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest4_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_5]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_5.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest5_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest5_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_6]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_6.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest6_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest6_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_7]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_7.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest7_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest7_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_8]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_8.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest8_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest8_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_9]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_9.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest9_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest9_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_10]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_10.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest10_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest10_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_11]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_11.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest11_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest11_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_12]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_12.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest12_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest12_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_13]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_13.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest13_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest13_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_14]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_14.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest14_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest14_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_15]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_15.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest15_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest15_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_16]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_16.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest16_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest16_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB

[program:postgrest_17]
command=/opt/botjagwar/bin/postgrest /opt/botjagwar/conf/pgrest/pgrest_17.ini
stdout_logfile=/opt/botjagwar/user_data/supervisor_postgrest17_stdout.log
stderr_logfile=/opt/botjagwar/user_data/supervisor_postgrest17_stderr.log
autostart=true
startretries=10
user=user
stderr_logfile_maxbytes=100MB